# elasticsearch5.4.0/osqueryd.conf
{
  "options": {
    "config_plugin": "filesystem",
    "logger_plugin": "filesystem",
    "utc": true
  },
  "schedule": {
   // This is a simple example query that outputs basic system information.
    "system_info": {
      // The exact query to run.
      "query": "SELECT hostname, cpu_brand, physical_memory FROM system_info;",
      // The interval in seconds to run this query, not an exact interval.
      "interval": 3600
    }
    ,"kernel_info": {
      "query" : "select * from kernel_info;",
      "interval": 300
    }
    ,"os_version": {
      "query" : "select * from os_version;",
      "interval": 300
    }
    ,"mounts": {
      "query" : "select * from mounts;",
      "interval": 300
    }
    ,"usb_devices": {
      "query" : "select * from usb_devices;",
      "interval": 300
    }
    ,"deb_packages": {
      "query" : "select * from deb_packages;",
      "interval": 300
    }
    ,"kernel_modules": {
      "query" : "select * from kernel_modules;",
      "interval": 300
    }
    ,"rpm_packages": {
      "query" : "select * from rpm_packages;",
      "interval": 300
    }
    ,"listening_ports": {
      "query" : "select * from listening_ports;",
      "interval": 300
    }
    ,"smbios_tables": {
      "query": "select * from smbios_tables;",
      "interval": 300
    }
    ,"pci_devices": {
      "query": "select * from pci_devices;",
      "interval": 300
    }
    ,"users": {
      "query": "select * from users;",
      "interval": 300
    }
    ,"groups": {
      "query": "select * from groups;",
      "interval": 300
    }
    ,"system_controls": {
      "query": "select * from system_controls;",
      "interval": 300
    }
    ,"sudoers": {
      "query": "select * from sudoers;",
      "interval": 300
    }
    ,"suid_bin": {
      "query": "select * from suid_bin;",
      "interval": 300
    }
    ,"startup_items": {
      "query": "select * from startup_items;",
      "interval": 300
    }
    ,"routes": {
      "query": "select * from routes;",
      "interval": 300
    }
    ,"logged_in_users": {
      "query": "select * from logged_in_users;",
      "interval": 300
    }
    ,"iptables": {
      "query": "select * from iptables;",
      "interval": 300
    }
    ,"block_devices": {
      "query": "select * from block_devices;",
      "interval": 300
    }
    ,"apt_sources": {
      "query": "select * from apt_sources;",
      "interval": 300
    }
  }
  ,"packs": {
    // "osquery-monitoring": "/usr/share/osquery/packs/osquery-monitoring.conf",
    // "incident-response": "/usr/share/osquery/packs/incident-response.conf",
    "it-compliance": "/usr/share/osquery/packs/it-compliance.conf"
    // "osx-attacks": "/usr/share/osquery/packs/osx-attacks.conf",
    // "vuln-management": "/usr/share/osquery/packs/vuln-management.conf",
    // "hardware-monitoring": "/usr/share/osquery/packs/hardware-monitoring.conf",
    // "ossec-rootkit": "/usr/share/osquery/packs/ossec-rootkit.conf"
  }
}
